# infra/systemd/afrolete-due-worker.service
[Unit]
Description=AfroLete due background worker
Documentation=/opt/afrolete/backend/README.md
After=network-online.target afrolete-backend.service
Wants=network-online.target

[Service]
Type=oneshot
User=afrolete
WorkingDirectory=/opt/afrolete/backend
EnvironmentFile=/run/pjs/afrolete-backend.env
Environment=AFROLETE_DUE_WORKER_LIMIT=25
Environment=AFROLETE_DUE_WORKER_COMMUNICATION_DIGEST_LIMIT=100
Environment=AFROLETE_DUE_WORKER_PERFORMANCE_LIMIT=25
Environment=AFROLETE_DUE_WORKER_PERFORMANCE_FORECAST_VALIDATION_LIMIT=25
Environment=AFROLETE_DUE_WORKER_PERFORMANCE_FORECAST_DRIFT_REPEAT_AFTER_HOURS=24
Environment=AFROLETE_DUE_WORKER_PERFORMANCE_FORECAST_DRIFT_CHANNEL=in_app
Environment=AFROLETE_DUE_WORKER_PERFORMANCE_REVIEW_LIMIT=25
Environment=AFROLETE_DUE_WORKER_PERFORMANCE_REVIEW_HORIZON_HOURS=24
Environment=AFROLETE_DUE_WORKER_PERFORMANCE_REVIEW_REPEAT_AFTER_HOURS=24
Environment=AFROLETE_DUE_WORKER_PERFORMANCE_INJURY_RISK_LIMIT=25
Environment=AFROLETE_DUE_WORKER_PERFORMANCE_INJURY_RISK_THRESHOLD=65
Environment=AFROLETE_DUE_WORKER_PERFORMANCE_INJURY_RISK_REPEAT_AFTER_HOURS=24
Environment=AFROLETE_DUE_WORKER_PERFORMANCE_INJURY_RISK_CHANNEL=in_app
Environment=AFROLETE_DUE_WORKER_WEARABLE_PULL_LIMIT=25
Environment=AFROLETE_DUE_WORKER_WEARABLE_PULL_MAX_PAGES=3
Environment=AFROLETE_DUE_WORKER_WEARABLE_PULL_DEFAULT_RETRY_AFTER_SECONDS=300
Environment=AFROLETE_DUE_WORKER_WEBHOOK_MAX_ATTEMPTS=5
ExecStart=/opt/afrolete/backend/.venv/bin/python -m app.workers.due \
  --limit ${AFROLETE_DUE_WORKER_LIMIT} \
  --communication-digest-limit ${AFROLETE_DUE_WORKER_COMMUNICATION_DIGEST_LIMIT} \
  --performance-limit ${AFROLETE_DUE_WORKER_PERFORMANCE_LIMIT} \
  --performance-forecast-validation-limit ${AFROLETE_DUE_WORKER_PERFORMANCE_FORECAST_VALIDATION_LIMIT} \
  --auto-alert-performance-forecast-drift \
  --performance-forecast-drift-repeat-after-hours ${AFROLETE_DUE_WORKER_PERFORMANCE_FORECAST_DRIFT_REPEAT_AFTER_HOURS} \
  --performance-forecast-drift-channel ${AFROLETE_DUE_WORKER_PERFORMANCE_FORECAST_DRIFT_CHANNEL} \
  --performance-review-limit ${AFROLETE_DUE_WORKER_PERFORMANCE_REVIEW_LIMIT} \
  --performance-review-horizon-hours ${AFROLETE_DUE_WORKER_PERFORMANCE_REVIEW_HORIZON_HOURS} \
  --performance-review-repeat-after-hours ${AFROLETE_DUE_WORKER_PERFORMANCE_REVIEW_REPEAT_AFTER_HOURS} \
  --performance-injury-risk-limit ${AFROLETE_DUE_WORKER_PERFORMANCE_INJURY_RISK_LIMIT} \
  --performance-injury-risk-threshold ${AFROLETE_DUE_WORKER_PERFORMANCE_INJURY_RISK_THRESHOLD} \
  --performance-injury-risk-repeat-after-hours ${AFROLETE_DUE_WORKER_PERFORMANCE_INJURY_RISK_REPEAT_AFTER_HOURS} \
  --performance-injury-risk-channel ${AFROLETE_DUE_WORKER_PERFORMANCE_INJURY_RISK_CHANNEL} \
  --wearable-pull-limit ${AFROLETE_DUE_WORKER_WEARABLE_PULL_LIMIT} \
  --wearable-pull-max-pages ${AFROLETE_DUE_WORKER_WEARABLE_PULL_MAX_PAGES} \
  --wearable-pull-default-retry-after-seconds ${AFROLETE_DUE_WORKER_WEARABLE_PULL_DEFAULT_RETRY_AFTER_SECONDS} \
  --webhook-max-attempts ${AFROLETE_DUE_WORKER_WEBHOOK_MAX_ATTEMPTS}
NoNewPrivileges=true
PrivateTmp=true
